# systemd/snapshot@.service
[Unit]
Description=OPSVI Git Snapshot Service for %i
After=network.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
WorkingDirectory=/home/opsvi/master_root
ExecStart=/home/opsvi/master_root/tools/snapshot.sh
Restart=no
RestartSec=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opsvi-snapshot

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/opsvi/master_root /home/opsvi/.opsvi_logs

[Install]
WantedBy=default.target